# deploy/owui-mcp-spawner.service.example
[Unit]
Description=OWUI MCP Spawner
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/owui-mcp-spawner

# Add --no-edit to block all code/config editing (read-only UI).
# Add --no-code-edit to allow JSON uploads but block the inline code editor.
ExecStart=/home/YOUR_USER/owui-mcp-spawner/.venv/bin/python app/manager.py --host 0.0.0.0

# Load password and other secrets from a separate file (not committed to git)
EnvironmentFile=/etc/owui-mcp-spawner.env

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
